# files/glusterfs.service
[Unit]
Description=GlusterFS Mount

[Service]
Type=forking
ExecStart=/bin/mount -a

[Install]
WantedBy=multi-user.target
