# files/keepalived.timer
[Unit]
Description=KeepAlived timer on boot

[Timer]
OnBootSec=20
Unit=keepalived.service

[Install]
WantedBy=multi-user.target
